# supervisor.conf
[program:bqm]
command=/home/ifnmg/www/bqm/start.sh
user=ifnmg
stdout_logfile=/home/ifnmg/www/bqm/logs/supervisor.log
redirect_stderr=true
autostart=true
autorestart=true
